# modes/ai/systemd/lmstudio.service
# Refract OS — LM Studio OpenAI-compatible server on port 8080.
#
# This is a USER service (LM Studio is per-user in ~/.lmstudio). Install:
#   mkdir -p ~/.config/systemd/user
#   cp lmstudio.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now lmstudio.service
#
# For an always-on headless box (server runs even when you're not logged in),
# enable lingering once:  sudo loginctl enable-linger "$USER"
#
# %h expands to the user's home dir. `lms server start` registers the server in
# the llmster daemon and returns, so Type=oneshot + RemainAfterExit is correct.

[Unit]
Description=Refract OS — LM Studio OpenAI server (port 8080)
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-%h/.lmstudio/bin/lms daemon up
ExecStart=%h/.lmstudio/bin/lms server start --port 8080
ExecStop=%h/.lmstudio/bin/lms server stop

[Install]
WantedBy=default.target
